# ubuntu_fs/var/lib/dpkg/info/util-linux.conffiles
/etc/init/hwclock.conf
/etc/init/hwclock-save.conf
/etc/cron.weekly/fstrim
